# collector-service/osquery/osquery.macos.conf
{
  "//": "EDRdog osquery 설정 (macOS). es_process_events(EndpointSecurity) + socket_events + file_events(FIM) 를 kafka_producer 로 events-raw 토픽에 전송한다.",
  "//kafka": "실제 Kafka 전송 플래그(--logger_plugin=kafka_producer 등)는 osquery.flags 로 넘긴다. 아래 kafka_topics 는 쿼리명->토픽 라우팅.",
  "//type": "collector 의 RawEventMapper 가 쿼리명으로 type 을 판정한다: socket/network->network, file->file, script->script, 그 외 process.",

  "options": {
    "disable_events": false,
    "disable_endpointsecurity": false,
    "host_identifier": "hostname",
    "schedule_splay_percent": 10
  },

  "kafka_topics": {
    "events-raw": ["process_events", "socket_events", "script_events", "file_events"]
  },

  "file_paths": {
    "autorun": [
      "/Users/%/Library/LaunchAgents/%%",
      "/Library/LaunchAgents/%%",
      "/Library/LaunchDaemons/%%"
    ]
  },

  "schedule": {
    "process_events": {
      "//": "es_process_events: 프로세스 생성 감시. parent 는 PID 라 processes 조인으로 이름을 넣어 collector 가 그대로 쓰게 한다.",
      "query": "SELECT e.path AS path, e.cmdline AS cmdline, p.name AS parent FROM es_process_events e LEFT JOIN processes p ON e.parent = p.pid WHERE e.event_type = 'exec';",
      "interval": 10,
      "removed": false
    },
    "socket_events": {
      "//": "아웃바운드 연결만(remote_address 존재, connect). C2 비컨 탐지 입력.",
      "query": "SELECT path, remote_address, remote_port FROM socket_events WHERE action = 'connect' AND remote_address != '' AND remote_address NOT IN ('127.0.0.1', '::1', '0.0.0.0');",
      "interval": 10,
      "removed": false
    },
    "script_events": {
      "//": "스크립트 인터프리터 실행. cmdline 에 스크립트 경로가 담겨, 임시/다운로드 경로 실행을 detector 가 MEDIUM(T1059)으로 판정한다. process_events 와 같은 evented 테이블이라 인터프리터 exec 는 양쪽에 뜰 수 있다(의도된 중복).",
      "query": "SELECT e.path AS path, e.cmdline AS cmdline, p.name AS parent FROM es_process_events e LEFT JOIN processes p ON e.parent = p.pid WHERE e.event_type = 'exec' AND (e.path LIKE '%/bash' OR e.path LIKE '%/sh' OR e.path LIKE '%/zsh' OR e.path LIKE '%/python%' OR e.path LIKE '%/osascript');",
      "interval": 10,
      "removed": false
    },
    "file_events": {
      "//": "FIM: 자동실행(LaunchAgents/Daemons) 경로의 파일 생성/변경 감시. target_path 전체를 detector 가 지속성(T1547) 판정에 쓴다. file_paths.autorun 참조.",
      "query": "SELECT target_path, action FROM file_events WHERE action IN ('CREATED', 'UPDATED', 'MOVED_TO');",
      "interval": 10,
      "removed": false
    }
  }
}
